# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:mysql]
command=/usr/bin/mysqld_safe --user=mysql
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/mysql.log
stderr_logfile=/var/log/supervisor/mysql.err
startsecs=5

[program:redis]
command=redis-server --daemonize no --loglevel notice
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis.err

[program:blazegraph]
command=java -server -Xmx2g -jar /opt/blazegraph/blazegraph.jar
directory=/opt/blazegraph
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/blazegraph.log
stderr_logfile=/var/log/supervisor/blazegraph.err
startsecs=10

[program:otnode]
command=/opt/ot-node/start-node.sh
directory=/opt/ot-node
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/supervisor/otnode.log
stderr_logfile=/var/log/supervisor/otnode.err
startsecs=60
